# network_scanner.service
[Unit]
Description=Network Scanner Web Application
After=network.target

[Service]
User=YOUR_USERNAME  # Replace with your username
Group=YOUR_GROUP    # Replace with your group
WorkingDirectory=/path/to/your/app  # Replace with your app directory
Environment="PATH=/path/to/your/venv/bin"  # Replace with your virtual environment path
ExecStart=/path/to/your/venv/bin/python wsgi.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 